# pi_agent/說明文件/systemd檔案/homepi-hls-www@.service
[Unit]
Description=HomePi HLS static web server (CORS) on :8088
After=network-online.target homepi-hls.service
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i

EnvironmentFile=/home/%i/pi_agent/.env
ExecStart=/usr/bin/python3 /home/%i/pi_agent/serve_hls.py
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target