# pi_agent/說明文件/systemd檔案/homepi-hls.service
[Unit]
Description=HomePi Camera HLS Stream (Low-Latency TS)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qkauia
WorkingDirectory=/home/qkauia/pi_agent/stream
EnvironmentFile=/home/qkauia/pi_agent/.env
ExecStartPre=/usr/bin/bash -lc 'mkdir -p /home/qkauia/pi_agent/stream && rm -f /home/qkauia/pi_agent/stream/index.m3u8.tmp /home/qkauia/pi_agent/stream/index.m3u8 /home/qkauia/pi>
ExecStart=/usr/bin/bash -lc '\
/usr/bin/rpicam-vid -t 0 \
  --width 1280 --height 720 --framerate 30 \
  --codec h264 --profile baseline --level 4.0 \
  --intra 15 --inline --nopreview \
  --libav-format h264 -o - \
| \
/usr/bin/ffmpeg -loglevel warning -fflags +genpts -use_wallclock_as_timestamps 1 \
  -f h264 -i - -c:v copy -an \
  -f hls -hls_time 0.5 -hls_list_size 3 -start_number 1 \
  -hls_flags delete_segments+append_list+program_date_time+independent_segments+temp_file \
  -hls_segment_filename /home/qkauia/pi_agent/stream/seg_%%05d.ts \
  /home/qkauia/pi_agent/stream/index.m3u8 \
'
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=10
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target
