# pi_agent/說明文件/systemd檔案/homepi-scheduler@.service
[Unit]
Description=HomePi Scheduler (pull schedules, run actions, auto-light) for user %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/pi_agent

EnvironmentFile=/home/%i/pi_agent/.env
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/python3 -m utils.scheduler

Restart=on-failure
RestartSec=2s
KillMode=control-group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target