# pi_agent/homepi-agent.service
[Unit]
Description=HomePi HTTP Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qkauia
WorkingDirectory=/home/qkauia/pi_agent

ExecStart=/usr/bin/python3 /home/qkauia/pi_agent/http_agent.py
Environment=PYTHONUNBUFFERED=1
Environment=GPIOZERO_PIN_FACTORY=lgpio
Environment=LED_PIN=17
Environment=LED_ACTIVE_HIGH=true

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
